# power-monitoring/power-monitoring.service
[Unit]
Description=Monitor power usage
After=network.target
StartLimitBurst=10
StartLimitIntervalSec=20

[Service]
Type=simple
Restart=always
RestartSec=1
User=root
ExecStart=/usr/bin/env bash /root/power-monitoring/manager.sh

[Install]
WantedBy=multi-user.target
